# docker/repo-update-tarball-install.cron
1 * * * * root flock -n /var/lock/repo-tarball-client /usr/bin/update_tarball-install.sh
@reboot   root flock -n /var/lock/repo-tarball-client /usr/bin/update_tarball-install.sh
