# docker/supervisor-rsync.conf
[program:rsync]
command=/usr/bin/rsync --daemon --no-detach
autorestart=true
